# de.inberlinwohnen.scraper.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
    "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>de.inberlinwohnen.scraper</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/Younes/Desktop/BerlinWohnen/.venv/bin/python</string>
        <string>/Users/Younes/Desktop/BerlinWohnen/scraper.py</string>
    </array>

    <key>StartInterval</key>
    <integer>300</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>WorkingDirectory</key>
    <string>/Users/Younes/Desktop/BerlinWohnen</string>

    <key>StandardOutPath</key>
    <string>/Users/Younes/Desktop/BerlinWohnen/launchd.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/Younes/Desktop/BerlinWohnen/launchd.log</string>
</dict>
</plist>
